# supervisord.conf
[supervisord]
nodaemon=true

[program:webadm]
command=/start.sh
autostart=true
autorestart=false

